# config/warden.service
[Unit]
Description=Warden Service


[Service]
ExecStart=/root/eth_node_client/warden.py